# data/fpemud-portal.service
[Unit]
Description=Web Portal for Fpemud
After=network.target

[Service]
ExecStart=/usr/sbin/fpemud-portal

[Install]
WantedBy=multi-user.target
